# scripts/automation/com.optionscalculator.evidence-watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.optionscalculator.evidence-watchdog</string>
  <key>ProgramArguments</key>
  <array>
    <string>__PROJECT_ROOT__/scripts/automation/run_daily_evidence_watchdog.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>22</integer>
    <key>Minute</key>
    <integer>15</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>__HOME__/.options_calculator_pro/logs/evidence_watchdog_launchd_stdout.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/.options_calculator_pro/logs/evidence_watchdog_launchd_stderr.log</string>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
